# debian/postinst
#!/bin/sh

systemctl daemon-reload
systemctl --global enable xfconnect.service

exit 0