# debian/postrm
#!/bin/sh

systemctl --global disable xfconnect.service
systemctl daemon-reload

exit 0