# source/systemd/user/xfconnect.service
[Unit]
Description=xfconnect indicator service
After=graphical.target
Requieres=xfce4-notifyd.service
Requires=dbus.service


[Service]
#Environment=PATH=uncomment:to:override:your:PATH
ExecStart=/usr/local/bin/xfconnect-indicator.py -s
TimeoutStartSec=2s
Restart=on-failure
RestartSec=3s

[Install]
WantedBy=default.target
